# root/etc/services.d/cameraproxy/run
#!/usr/bin/with-contenv bash

umask 022

cd /config || exit

exec \
	s6-setuidgid abc mono --debug MJpegCameraProxyCmd.exe
